# home-assistant-configurator@pi.service.in
[Unit]
Description=Home Assistant
After=network.target

[Service]
Type=simple
User=%i
ExecStart=<path to homeassistant configuration directory>/configurator-start.sh

[Install]
WantedBy=multi-user.target
