# home-assistant-appdaemon@pi.service.in
[Unit]
Description=AppDaemon
After=home-assistant@pi.service

[Service]
Type=simple
User=%I
ExecStart=<path to homeassistant configuration directory>/appdaemon/.venv/bin/appdaemon -c <path to homeassistant configuration directory>/appdaemon
Restart=on-failure

[Install]
WantedBy=multi-user.target
