# home-assistant@pi.service.in
[Unit]
Description=Home Assistant
After=network.target

[Service]
Type=simple
User=%i
ExecStart=<path to homeassistant configuration directory>/start.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
